# mica-dist/src/main/deb/debian/postinst
#!/bin/sh
# postinst script for mica2
#
# see: dh_installdeb(1)

set -e

. /usr/share/debconf/confmodule

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
    configure)

    	RET=""
      db_get mica-server/admin_password

      if [ -z "$2" ] && [ "$RET" != "" ] ; then
        echo Updating administrator password.
        # hash the password and escape the possible '$/&' characters
        adminpw=$(echo -n $RET | xargs java -jar /usr/share/mica2-*/tools/lib/obiba-password-hasher-*-cli.jar | sed -e 's/\([\$\/\/]\)/\\\1/g' | cut -d ' ' -f 1)
        tfile=`mktemp`
        if [ ! -f "$tfile" ]; then
          return 1
        fi
        cat /etc/mica2/shiro.ini | sed -e "s,^administrator\s*=.*\,,administrator=$adminpw\,," > $tfile && \
          mv $tfile /etc/mica2/shiro.ini
        # clear password in db_get so it won't be re-encoded during next upgrade
        db_set mica-server/admin_password ""
      fi

    	# Create mica user if it doesn't exist.
      if ! id mica > /dev/null 2>&1 ; then
        adduser --system --home /var/lib/mica2 --no-create-home --disabled-password mica
      fi

      # Mica Server file structure on Debian
      # /etc/mica2: configuration
      # /usr/share/mica2: executable
      # /var/lib/mica2: data runtime
      # /var/log/mica2: logs

      rm -f /usr/share/mica2
      new_release="$(ls -t /usr/share/ |grep mica2|head -1)"
      ln -s /usr/share/${new_release} /usr/share/mica2

      if [ ! -e /var/lib/mica2/conf ] ; then
        ln -s /etc/mica2 /var/lib/mica2/conf
      fi

      # Some deprecated configuration files
      if [ -e /var/lib/mica2/conf/mica-categories.yml ] ; then
        mv -f /var/lib/mica2/conf/mica-categories.yml /var/lib/mica2/conf/mica-categories.yml.old
      fi

      if [ -e /var/lib/mica2/conf/mica-aggregations.yml ] ; then
        mv -f /var/lib/mica2/conf/mica-aggregations.yml /var/lib/mica2/conf/mica-aggregations.yml.old
      fi

      chown -R mica:adm /var/lib/mica2 /var/log/mica2 /etc/mica2
      chmod -R 750      /var/lib/mica2 /var/log/mica2 /etc/mica2
      chmod +x /usr/share/mica2/tools/shiro-hasher
      find /etc/mica2 -type f -print0 | xargs -0 chmod 640

      if [ -e /etc/init.d/mica2 ]; then
        update-rc.d mica2 defaults >/dev/null
        rm /etc/init.d/mica2
      fi
      systemctl daemon-reload
      systemctl enable mica2
      systemctl start mica2
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
